# crontab
1/10 * * * * kill -52 $(pidof dwmblocks); randwp
1/30 * * * * kill -49 $(pidof dwmblocks)
